# escenario4/maquina1/kafka-zookeeper.service
[Unit]
Description=Apache Zookeeper server (Kafka)
Documentation=https://zookeeper.apache.org
Requires=network.target remote-fs.target
After=network.target remote-fs.target

[Service]
Type=simple
ExecStart=/usr/local/src/kafka/bin/zookeeper-server-start.sh /usr/local/src/kafka/config/zookeeper.properties
ExecStop=/usr/local/src/kafka/bin/zookeeper-server-stop.sh

[Install]
WantedBy=multi-user.target
